# root/etc/services.d/httpd/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc /webui --port 80
