# root/etc/services.d/darkhttpd/finish
#!/usr/bin/execlineb -S0

exec \
    s6-svscanctl -t /var/run/s6/services